# debian/swift-container.swift-container-reconciler.upstart.in
# swift-container-reconciler - SWIFT Container Replicator
#
# The swift container reconciler.

description     "SWIFT Container Reconciler"
author          "Marc Cluet <marc.cluet@ubuntu.com>"

start on runlevel [2345]
stop on runlevel [016]

pre-start script
  if [ -f "/etc/swift/container-reconciler.conf" ]; then
    exec /usr/bin/swift-init container-reconciler start
  else
    exit 1
  fi
end script

post-stop exec /usr/bin/swift-init container-reconciler stop
