# debian/swift-object-expirer.init.in
#! /bin/sh
### BEGIN INIT INFO
# Provides:          swift-object-expirer
# Required-Start:    $remote_fs
# Required-Stop:     $remote_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Swift object expirer
# Description:       Object expirer daemon for Swift.
### END INIT INFO

PROJECT_NAME=swift
NAME="${PROJECT_NAME}-object-expirer"
DESC="Swift object expirer"
DAEMON_ARGS=/etc/${PROJECT_NAME}/object-expirer.conf
NO_OPENSTACK_CONFIG_FILE_DAEMON_ARG=yes
USE_LOGFILE=no
NO_OPENSTACK_LOGFILE_DAEMON_ARG=yes
