# debian/swift-object.swift-object-updater.init.in
#! /bin/sh
### BEGIN INIT INFO
# Provides:          swift-object-updater
# Required-Start:    $remote_fs
# Required-Stop:     $remote_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Swift object updater
# Description:       Object updater server for Swift.
### END INIT INFO

PROJECT_NAME=swift
NAME="${PROJECT_NAME}-object-updater"
DESC="Swift object updater"
DAEMON_ARGS=/etc/${PROJECT_NAME}/object-server.conf
NO_OPENSTACK_CONFIG_FILE_DAEMON_ARG=yes
USE_LOGFILE=no
NO_OPENSTACK_LOGFILE_DAEMON_ARG=yes
