# debian/swift.postrm
#!/bin/sh

set -e

invoke-rc.d rsyslog restart

#DEBHELPER#
