# debian/swift-object.swift-object-reconstructor.init.in
#! /bin/sh
### BEGIN INIT INFO
# Provides:          swift-object-reconstructor
# Required-Start:    $remote_fs
# Required-Stop:     $remote_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Swift object reconstructor
# Description:       Object reconstructor server for Swift.
### END INIT INFO

PROJECT_NAME=swift
NAME="${PROJECT_NAME}-object-reconstructor"
DESC="OpenStack Swift object reconstructor"
CONFIG=object-server
MULTI=1
