# debian/swift-object.init.in
#! /bin/sh
### BEGIN INIT INFO
# Provides:          swift-object
# Required-Start:    $remote_fs
# Required-Stop:     $remote_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Swift object server
# Description:       Object server for Swift.
### END INIT INFO

PROJECT_NAME=swift
NAME="${PROJECT_NAME}-object-server"
DESC="Swift object server"
SCRIPTNAME=/etc/init.d/swift-object
DAEMON_ARGS=/etc/${PROJECT_NAME}/object-server.conf
NO_OPENSTACK_CONFIG_FILE_DAEMON_ARG=yes
USE_LOGFILE=no
