# debian/swift.cron.d
*/5 * * * *     swift   test -x /usr/bin/swift-recon-cron && /usr/bin/swift-recon-cron /etc/swift/object-server.conf