# debian/swift-object.swift-object-updater.init.in
#! /bin/sh
### BEGIN INIT INFO
# Provides:          swift-object-updater
# Required-Start:    $remote_fs
# Required-Stop:     $remote_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Swift object updater
# Description:       Object updater server for Swift.
### END INIT INFO

PROJECT_NAME=swift
NAME="${PROJECT_NAME}-object-updater"
DESC="OpenStack Swift object updater"
CONFIG=object-server
MULTI=1
SYSTEMD_DOC="man:${NAME}(1)"
